# config/chroot/etc/systemd/system/boot-config.service
[Service]
Type=idle
ExecStart=/usr/local/sbin/boot-config

[Install]
WantedBy=multi-user.target
